# systemd/offlineimap.timer
[Timer]
OnUnitInactiveSec=60s
Unit=offlineimap.service
